# deploy/dearyou-scheduler.timer
[Unit]
Description=Run the DearYou Laravel scheduler every minute

[Timer]
OnCalendar=*-*-* *:*:00
Persistent=true
Unit=dearyou-scheduler.service

[Install]
WantedBy=timers.target
